# scripts/automation/com.memory-router.source-distillation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.memory-router.source-distillation</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{BUN_PATH}}</string>
        <string>run</string>
        <string>src/cli/distill-sources.ts</string>
        <string>--apply</string>
    </array>
    <key>WorkingDirectory</key>
    <string>{{PROJECT_ROOT}}</string>
    <key>StartInterval</key>
    <integer>{{INTERVAL_SECONDS}}</integer>
    <key>StandardOutPath</key>
    <string>{{PROJECT_ROOT}}/logs/source-distillation.log</string>
    <key>StandardErrorPath</key>
    <string>{{PROJECT_ROOT}}/logs/source-distillation.log</string>
    <key>RunAtLoad</key>
    <false/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
